# jobs/boshhmforwarder/templates/boshhmforwarder_logrotate.cron.erb
*/5 * * * * test -x /usr/sbin/logrotate && /usr/sbin/logrotate /var/vcap/jobs/boshhmforwarder/config/logrotate.conf >>/var/vcap/sys/log/boshhmforwarder/boshhmforwarder_logrotate_cron.log 2>&1
